# root/etc/s6-overlay/s6-rc.d/init-mod-calibre-web-tolino/up
sed 's/^[^#]*lm\.session_protection/#&/' -i /app/calibre-web/cps/__init__.py